# srv-os3-ost-6.mount
[Unit]
Description=Mount MDTs and OSTs
After=srv-os3-mdt.mount

[Mount]
What=/dev/sdg
Where=/srv/os3/ost/6
Type=lustre
Options=defaults
