# srv-os3-ost-5.mount
[Unit]
Description=Mount MDTs and OSTs
After=srv-os3-mdt.mount

[Mount]
What=/dev/sdf
Where=/srv/os3/ost/5
Type=lustre
Options=defaults
