# srv-os3-ost-4.mount
[Unit]
Description=Mount MDTs and OSTs
After=lustre.service

[Mount]
What=/dev/sde
Where=/srv/os3/ost/4
Type=lustre
Options=defaults
