# srv-os3-ost-3.mount
[Unit]
Description=Mount MDTs and OSTs
After=lustre.service

[Mount]
What=/dev/sdd
Where=/srv/os3/ost/3
Type=lustre
Options=defaults
