# etc/systemd/srv-os3-ost-4.mount
[Unit]
Description=Mount MDTs and OSTs
After=srv-os3-mdt.mount

[Mount]
What=/dev/sde
Where=/srv/os3/ost/4
Type=lustre
Options=defaults
