# etc/systemd/srv-os3-ost-2.mount
[Unit]
Description=Mount MDTs and OSTs
After=srv-os3-mdt.mount

[Mount]
What=/dev/sdc
Where=/srv/os3/ost/2
Type=lustre
Options=defaults
